# docker/sakuli-omd-labs/omd_ansible_drop-in_role/files/skel/etc/cron.d/sakuli_gearman_proxy
MAILTO=""
*/1 * * * * perl $OMD_ROOT/local/bin/appradar_gearman_proxy.pl -c $OMD_ROOT/etc/mod-gearman/appradar_gearman_proxy.cfg &

